# assets/training/model_management/environments/mlflow-model-inference/context/runit/rsyslog/run
#!/bin/bash

echo "`date -uIns` - rsyslog/run $@"

exec "${AML_LOGGER_ROOT:-/var/azureml-logger}/start_logger.sh"